# tortuga_kits/awsadapter_6_3_0/puppet_modules/univa-tortuga_kit_awsadapter/templates/awsspotd.service.erb
[Unit]
Description=AWS resource adapter spot instance support daemon
After=syslog.target network.target auditd.service

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/awsspotd
PIDFile=/var/run/awsspotd.pid
ExecStart=<%= @instroot -%>/bin/awsspotd $AWSSPOTD_OPTIONS --pidfile /var/run/awsspotd.pid
KillMode=process

[Install]
WantedBy=multi-user.target
